# daemon/web_sheets_storage.conf
description "web_sheets storage server"
author      "Charles Rymal"

start on runlevel [2345]
stop on shutdown

script
    exec web_sheets_storage
end script

pre-start script
end script

pre-stop script
end script


