# daemon/web_sheets_storage.service
[Unit]
Description=web_sheets storage server

[Service]
ExecStart=/usr/local/bin/web_sheets_storage


