# daemon/web_sheets_sheets_backend.service
[Unit]
Description=web_sheets sheets_backend server

[Service]
User=web_sheets
ExecStart=/usr/local/bin/web_sheets_sheets_backend


